# overlay/etc/systemd/system/datadir.mount
[Unit]
Before=local-fs.target

[Mount]
What=tmpfs
Where=/datadir
Type=tmpfs
Options=defaults,noatime,size=50g,mode=1777
